# packaging/worker/debian/postinst
#!/bin/bash
set -e

# Create system user if it doesn't exist
if ! id -u utworker >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin utworker
fi

# Create config directory
mkdir -p /etc/ultimate-terminal
chmod 755 /etc/ultimate-terminal

# Create default config if not exists
if [ ! -f /etc/ultimate-terminal/worker.env ]; then
    cat > /etc/ultimate-terminal/worker.env << 'EOF'
# Ultimate Terminal Worker Configuration
# Edit this file and restart the service

# Nexus server URL (required)
NEXUS_URL=http://localhost:3002

# Worker authentication key (Generated in Nexus)
API_KEY=

# Heartbeat interval in milliseconds
# WORKER_HEARTBEAT_MS=5000

# Shell to use (defaults to user's shell or bash)
# SHELL=/bin/bash
EOF
    chmod 600 /etc/ultimate-terminal/worker.env
fi

# Enable and start service (only if systemd is available)
if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload
  systemctl enable ultimate-terminal-worker.service || true
fi

echo ""
echo "╔══════════════════════════════════════════════════════════════╗"
echo "║      Ultimate Terminal Worker instalado correctamente       ║"
echo "╠══════════════════════════════════════════════════════════════╣"
echo "║  1. Configura tu API Key:                                   ║"
echo "║     sudo nano /etc/ultimate-terminal/worker.env             ║"
echo "║                                                             ║"
echo "║  2. Inicia el servicio:                                     ║"
echo "║     sudo systemctl start ultimate-terminal-worker           ║"
echo "║                                                             ║"
echo "║  3. Verifica el estado:                                     ║"
echo "║     sudo systemctl status ultimate-terminal-worker          ║"
echo "║                                                             ║"
echo "║  4. Ver logs:                                               ║"
echo "║     sudo journalctl -u ultimate-terminal-worker -f          ║"
echo "╚══════════════════════════════════════════════════════════════╝"
echo ""

exit 0
